# examples/ROOT/img/tmpfs/service/getty/run
#!/command/execlineb -P
getty -L -H DEBUG 38400 tty1
